# BotenODK-ObjectDetectionKit/crontab
*/5 * * * * python /usr/src/app/process_image.py